# backup-repo-template/net.xzer.work-backup-hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.xzer.work-backup-hourly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/work-machine-backup/scripts/backup-runner.sh</string>
        <string>/path/to/work-backup</string>
        <string>--commit-only</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
